# base_rootfs/etc/inittab
::sysinit:/etc/init.d/rcS

::respawn:-/bin/sh
#tty2::askfirst:-/bin/sh
#tty3::askfirst:-/bin/sh

::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a

