# base_rootfs/etc/init.d/rc.local
#!/bin/sh

telnetd
mkdir /var/log/boa
boa